# home/dotfiles/vim/deleteundofiles.timer
[Unit]
Description="Timer for deleteundofiles.service"

[Timer]
Unit=deleteundofiles.service
OnUnitActiveSec=1day

[Install]
WantedBy=timers.target
